# crontab
* * * * *  /usr/bin/php /opt/support-ticket/artisan schedule:run >> /dev/null 2>&1
